# docker/supervisord.conf
[program:atlas]
command=node index.js
directory=/app
autostart=true
autorestart=true
environment=NODE_ENV=production
stderr_logfile=/var/log/atlas.err.log
stdout_logfile=/var/log/atlas.out.log
